# supervisord.conf
[supervisord]
nodaemon=true

[program:dockapp]
command=./dockapp-start

[program:mediacenter]
command=/bin/bash -c "mediacenter21 /mediaserver"
